# supervisor.ini
[program:picostats]
command=/picostats
stderr_logfile=/picostats.err
stdout_logfile=/picostats.log
directory=/

[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/redis/stdout.log
stderr_logfile=/var/log/redis/stderr.log

[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgresql/data -c config_file=/var/lib/postgresql/data/postgresql.conf
user=postgres
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
username=www-data
autorestart=true
autostart=true
